# etc/systemd_overrides/dnscrypt-proxy.socket.d/override.conf
[Socket]
ListenStream=
ListenDatagram=
#ListenStream=127.0.0.1:40
#ListenDatagram=127.0.0.1:40
ListenStream=127.0.0.1:53
ListenDatagram=127.0.0.1:53
